# container-nordvpn.service
# container-nordvpn.service
# autogenerated by Podman 5.2.3
# Mon May 26 10:31:48 BST 2025

[Unit]
Description=Podman container-nordvpn.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/nix/store/4j9b9m9frv8scsgksbhyibrhmbanxghr-podman-5.2.3/bin/.podman-wrapped start nordvpn
ExecStop=/nix/store/4j9b9m9frv8scsgksbhyibrhmbanxghr-podman-5.2.3/bin/.podman-wrapped stop  \
	-t 10 nordvpn
ExecStopPost=/nix/store/4j9b9m9frv8scsgksbhyibrhmbanxghr-podman-5.2.3/bin/.podman-wrapped stop  \
	-t 10 nordvpn
PIDFile=/run/containers/storage/overlay-containers/7cf853ca735af5e028db0da64f27cd75a13c372bf321c58bec5a74547a59de26/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
